# tarot-proxy.service
[Unit]
Description=Tarot Proxy Server
After=network.target

[Service]
Type=simple
User=proxy
WorkingDirectory=/opt/tarot-proxy-server
ExecStart=/opt/tarot-proxy-server/tarot_proxy_server
Restart=on-failure
EnvironmentFile=/opt/tarot-proxy-server/.env
StandardOutput=file:/var/log/tarot-proxy.log
StandardError=inherit

[Install]
WantedBy=multi-user.target